# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable occupancy-exporter.service
systemctl restart occupancy-exporter.service
