# packaging/scripts/preremove.sh
systemctl stop occupancy-exporter.service
systemctl disable occupancy-exporter.service
